# units/jupyterhub.service
[Unit]
Description=JupyterHub

[Service]
User=root
ExecStart=/io/repo/units/jupyterhub.sh -f /io/repo/jupyterhub_config.py
WorkingDirectory=/tmp
